# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
autorestart=true
priority=10

[program:traffic_app]
command=/app/start.sh
directory=/app
user=appuser
stdout_logfile=/var/log/supervisor/traffic_app.log
stderr_logfile=/var/log/supervisor/traffic_app.log
autorestart=true
priority=20
environment=PYTHONPATH="/app"

[program:data_updater]
command=python3 scripts/generate_web_assets.py
directory=/app
user=appuser
stdout_logfile=/var/log/supervisor/data_updater.log
stderr_logfile=/var/log/supervisor/data_updater.log
autorestart=false
priority=30
environment=PYTHONPATH="/app"